# bootstrap/eink-bootstrap.service
[Unit]
Description=E-ink Display Bootstrapper
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root
Environment=PYTHONPATH=/opt/eink/venv/lib/python3.9/site-packages
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/opt/eink
ExecStart=/bin/bash /opt/eink/bootstrap/bootstrap.sh
ExecStop=/bin/true
# Run every hour for updates and health checks
ExecStartPost=/bin/bash -c 'systemd-run --on-calendar="*:0/1" --unit=eink-bootstrap-check.service /bin/bash /opt/eink/bootstrap/bootstrap.sh'
# Allow manual configuration trigger
ExecStartPost=/bin/bash -c 'systemctl enable eink-bootstrap-check.service'

# Security settings
ProtectSystem=full
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target 